# sv/vmware/run
#!/bin/sh
/etc/init.d/vmware start
exec chpst -b vmware pause
