# sv/vmware-networks/finish
#!/bin/sh
/usr/bin/vmware-networks --stop
