# sv/vmware-usbarbitrator/run
#!/bin/sh
/usr/lib/vmware/bin/vmware-usbarbitrator -f
